# grade-calculator.service
[Unit]
Description=gradecalc
After=network.target

[Service]
ExecStart=/home/kaiz/startgradecalc.sh

[Install]
WantedBy=multi-user.target